# camrotate.path
# SJP  26 September, 2016
[Unit]
Description=Monitor webcam landing directory

[Path]
PathModified=/home/st33v/farm/cam/

[Install]
WantedBy=multi-user.target

